# dexter-log-stream.service
[Unit]
Description=Dexter AI Log Stream Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/dexter-ai
ExecStart=/usr/bin/node /opt/dexter-ai/log-stream-server.js
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3003

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/dexter-ai

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dexter-log-stream

[Install]
WantedBy=multi-user.target