# dexter-log-api.service
[Unit]
Description=Dexter Log API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/dexter
ExecStart=/usr/bin/node /opt/dexter/log-api-server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dexter-log-api

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/dexter
ReadOnlyPaths=/opt/dexter

# Resource limits
MemoryMax=512M
TasksMax=100

# Environment
Environment=NODE_ENV=production
Environment=PORT=3004
Environment=LOG_DIR=/var/log/dexter

[Install]
WantedBy=multi-user.target